# sys/ARCity-auto-git-pull.service
[Unit]
Description=Auto Git Pull Service
After=network.target

[Service]
Type=simple
ExecStart=/opt/ARCity/auto-pull.sh
WorkingDirectory=/opt/ARCity/